# supervisord.conf
[supervisord]
nodaemon=true

[program:lambda1]
command=/usr/local/bin/aws-lambda-rie python -m awslambdaric happyRobot.lambda_function.lambda_handler
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr


[program:dashboard]
command=streamlit run dashboard/app.py --server.port=8501 --server.address=0.0.0.0
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
